# assets/etc/profile.d/apt-update.sh

#
# If there is no package cache or
# the package cache is older than one day,
# and a sources.list is present,
# update the package list.
#
# TODO
