# assets/etc/profile.d/path.sh

#
# Also look for executables in the current user's .local/bin
#
export PATH=$PATH:/home/$(whoami)/.local/bin
