# assets/etc/profile.d/ssh-keygen.sh

#
# If no SSH key pair is present,
# but the openssh-client is installed,
# generate an SSH key pair.
#
# TODO
