# assets/etc/profile.d/ssh-host-keygen.sh

#
# If the current user is root and
# the openssh-server is installed,
# but the host has no key pair,
# generate an SSH host key pair.
#
# TODO
